# slopd.service
[Unit]
Description=slopd — Claude agent session manager
Documentation=https://github.com/futpib/slopd
After=default.target

[Service]
ExecStart=%h/.cargo/bin/slopd
Restart=on-failure
RestartSec=5s
KillMode=process

[Install]
WantedBy=default.target
